# xd-torrent/pkg/xd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/XD"
daemon_user="_xd"
daemon_flags="${SYSCONFDIR}/xd/torrents.ini"

. /etc/rc.d/rc.subr

rc_bg=NO
rc_reload=NO
rc_stop_signal=KILL

rc_cmd $1